# street-server.service
[Unit]
Description=Street Asset Detection Inference Server
After=network.target network-online.target

[Service]
Type=simple
User=ankaa
Group=ankaa
WorkingDirectory=/home/ankaa/street/server
ExecStart=/bin/bash /home/ankaa/street/server/start.sh
Restart=always
RestartSec=10
Environment="PYTHONUNBUFFERED=1"

# Helpful limit increases for AI workloads
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
